# supervisor/hookbot.conf
[program:hookbot]
command=nohup manage.py hookbot &
autostart=true
autorestart=true
stderr_logfile=/var/log/hookbot.err.log
stdout_logfile=/var/log/hookbot.out.log