# supervisor/poolbot.conf
[program:poolbot]
command=nohup manage.py poolbot &
autostart=true
autorestart=true
stderr_logfile=/var/log/poolbot.err.log
stdout_logfile=/var/log/poolbot.out.log